# dozzle/rootfs/etc/services.d/dozzle/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App — Dozzle
# Ingress  : nginx on :8080 proxies here (:8081), patches HTML for SPA compat
# Agent    : optional built-in dozzle agent (enable_agent option)
# s6-overlay: https://github.com/just-containers/s6-overlay
# ==============================================================================

set -e

declare -a DOZZLE_ARGS
declare DOZZLE_ADDR
declare LOG_LEVEL
declare FILTER
declare NO_ANALYTICS
declare ENABLE_ACTIONS
declare INGRESS_ENTRY
declare REMOTE_AGENTS
declare AGENT_PORT
declare AGENT_HOSTNAME
declare ENABLE_AGENT

# --- Docker socket (docker_api: true in config.yaml) -------------------------
if [ -S /var/run/docker.sock ]; then
    if curl -s --unix-socket /var/run/docker.sock http://localhost/version >/dev/null 2>&1; then
        DOCKER_VERSION="$(curl -s --unix-socket /var/run/docker.sock http://localhost/version | jq -r '.Version')"
        bashio::log.info "Docker API OK — engine ${DOCKER_VERSION}"
    else
        bashio::log.warning "Docker socket found but API is not responding"
    fi
else
    bashio::log.error "No socket at /var/run/docker.sock — enable docker_api in the manifest"
fi

# --- Ingress entry (external URL prefix, stripped by Supervisor before forwarding) --
INGRESS_ENTRY="$(bashio::addon.ingress_entry)"

# nginx on :8080 → Dozzle on :8081 (internal only)
DOZZLE_ADDR=":8081"

# --- Read options from HA config ---------------------------------------------
LOG_LEVEL="$(bashio::config 'log_level' 'info')"
FILTER="$(bashio::config 'filter' '')"
NO_ANALYTICS="$(bashio::config 'no_analytics' 'true')"
ENABLE_ACTIONS="$(bashio::config 'enable_actions' 'false')"
ENABLE_AGENT="$(bashio::config 'enable_agent' 'false')"
AGENT_PORT="$(bashio::config 'agent_port' '7007')"
AGENT_HOSTNAME="$(bashio::config 'agent_hostname' '')"
REMOTE_AGENTS="$(bashio::config 'remote_agents' '')"

# --- Build Dozzle argument list ----------------------------------------------
DOZZLE_ARGS=(
    "--addr"  "${DOZZLE_ADDR}"
    "--base"  "/"
    "--level" "${LOG_LEVEL}"
    "--auth-provider" "none"
)

if bashio::var.true "${NO_ANALYTICS}"; then
    DOZZLE_ARGS+=("--no-analytics")
fi

if bashio::var.has_value "${FILTER}"; then
    DOZZLE_ARGS+=("--filter" "${FILTER}")
fi

if bashio::var.true "${ENABLE_ACTIONS}"; then
    DOZZLE_ARGS+=("--enable-actions")
fi

# --- Export env vars used by Dozzle ------------------------------------------
export DOZZLE_BASE="/"
export DOZZLE_ADDR="${DOZZLE_ADDR}"
export DOZZLE_LEVEL="${LOG_LEVEL}"
export DOZZLE_AUTH_PROVIDER="none"

if bashio::var.true "${NO_ANALYTICS}"; then
    export DOZZLE_NO_ANALYTICS="true"
fi

if bashio::var.has_value "${REMOTE_AGENTS}"; then
    export DOZZLE_REMOTE_AGENT="${REMOTE_AGENTS}"
fi

# --- Built-in agent (optional) -----------------------------------------------
if bashio::var.true "${ENABLE_AGENT}"; then
    bashio::log.info "Starting built-in Dozzle agent on 0.0.0.0:${AGENT_PORT} ..."
    if bashio::var.has_value "${AGENT_HOSTNAME}"; then
        DOZZLE_HOSTNAME="${AGENT_HOSTNAME}" \
            /usr/local/bin/dozzle agent \
                --addr "0.0.0.0:${AGENT_PORT}" \
                --level "${LOG_LEVEL}" &
    else
        /usr/local/bin/dozzle agent \
            --addr "0.0.0.0:${AGENT_PORT}" \
            --level "${LOG_LEVEL}" &
    fi
    sleep 1
    bashio::log.info "Built-in agent ready on port ${AGENT_PORT}"
fi

# --- Startup banner ----------------------------------------------------------
bashio::log.info "================================================"
bashio::log.info "  Dozzle for Home Assistant"
bashio::log.info "  App version  : ${BUILD_VERSION:-unknown}"
bashio::log.info "  Dozzle binary: ${DOZZLE_VERSION:-unknown}"
bashio::log.info "------------------------------------------------"
bashio::log.info "  Ingress (ext): ${INGRESS_ENTRY}"
bashio::log.info "  nginx proxy  : :8080  →  Dozzle :8081"
bashio::log.info "  Log level    : ${LOG_LEVEL}"
bashio::log.info "  Filter       : ${FILTER:-none}"
bashio::log.info "  Actions      : ${ENABLE_ACTIONS}"
bashio::log.info "  Agent        : ${ENABLE_AGENT} (port ${AGENT_PORT})"
bashio::log.info "  Remote agents: ${REMOTE_AGENTS:-none}"
bashio::log.info "================================================"
bashio::log.info "CMD: /usr/local/bin/dozzle ${DOZZLE_ARGS[*]}"

exec /usr/local/bin/dozzle "${DOZZLE_ARGS[@]}"
